# deploy/econdelta-aggregate.service.d/10-claude-json-writable.conf
[Service]
# 2026-05-29: same ~/.claude.json EROFS fix as parse — Opus review needs to write it.
ReadWritePaths=/home/adnan-local/.claude.json
