# deploy/econdelta-auction.service.d/10-claude-json-writable.conf
[Service]
# bb_auction calls the claude CLI two ways — the haiku CAPTCHA solver (via
# scrapers.bb_forex) and run_max (claude_max.max_client) for LLM row extraction.
# The CLI writes ~/.claude.json each run; ProtectHome=read-only blocks it (EROFS)
# and aborts the scrape. The .claude/ dir is covered by the main unit's
# ReadWritePaths; this adds the sibling .json file. Same carve-out as parse/
# aggregate/briefing — see AGENT_LEARNINGS 2026-05-29.
ReadWritePaths=/home/adnan-local/.claude.json
