# deploy/econdelta-parse.service.d/10-claude-json-writable.conf
[Service]
# 2026-05-29: claude CLI writes ~/.claude.json each run; ProtectHome=read-only blocks it
# (EROFS) and fails the parse preflight. The .claude/ dir was already carved out; this adds
# the sibling .json file. See AGENT_LEARNINGS 2026-05-29.
ReadWritePaths=/home/adnan-local/.claude.json
