# deploy/econdelta-media-screen.service.d/10-claude-json-writable.conf
[Service]
# media_screen calls run_max (claude_max.max_client) for per-article LLM
# extraction (media_screen/extract.py). The claude CLI writes ~/.claude.json
# each run; ProtectHome=read-only blocks it (EROFS) and aborts the screen
# before any extraction happens. The .claude/ dir is covered by the main unit's
# ReadWritePaths target; this adds the sibling .json file. Same carve-out as
# parse/aggregate/briefing/auction — see AGENTS landmine 17, AGENT_LEARNINGS 2026-05-29.
ReadWritePaths=/home/adnan-local/.claude.json
