# deploy/econdelta-parse.service
[Unit]
Description=EconDelta — Stage 2 parse (deterministic + Sonnet hybrid)
Documentation=https://github.com/clauding-lab/econdelta
After=econdelta-fetch.service network-online.target
Wants=network-online.target

# Retry on transient failures (Sonnet API hiccups, network blips). Up to 3
# starts per 30-min window — initial fire plus 2 retries. Note: parse needs
# fetched artifacts; if fetch failed, parse will run on the prior day's cache.
StartLimitIntervalSec=1800
StartLimitBurst=3

[Service]
Type=oneshot
User=adnan
Group=adnan
WorkingDirectory=/home/adnan/econdelta
EnvironmentFile=/etc/econdelta.env
# CLAUDE_BINARY lets parse_all find `claude` from cron/systemd PATH
Environment="CLAUDE_BINARY=/home/adnan/.npm-global/bin/claude"
ExecStart=/home/adnan/econdelta/.venv/bin/python -m parse_all

# Resource limits — pdfplumber + claude -p subprocess
MemoryMax=1500M
MemoryHigh=1200M
CPUQuota=200%
TasksMax=512

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/adnan/econdelta/data /home/adnan/econdelta/logs

# Timeouts + retry on transient failures
TimeoutStartSec=1800
Restart=on-failure
RestartSec=300

# Logging
StandardOutput=append:/home/adnan/econdelta/logs/parse-systemd.log
StandardError=append:/home/adnan/econdelta/logs/parse-systemd.log
SyslogIdentifier=econdelta-parse

[Install]
WantedBy=multi-user.target
