# deploy/econdelta-dse-dayend.timer
[Unit]
Description=EconDelta DSE DS30 per-ticker close scraper — daily 23:20 UTC (05:20 BDT next day), after the 23:11 index scrape
Requires=econdelta-dse-dayend.service

[Timer]
OnCalendar=*-*-* 23:20:00 UTC
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
