# deploy/econdelta-aggregate.timer
[Unit]
Description=Run EconDelta aggregator daily at 23:20 UTC (05:20 BDT next day)
Requires=econdelta-aggregate.service

[Timer]
OnCalendar=*-*-* 23:20:00 UTC
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
