# deploy/econdelta-media-screen.service
[Unit]
Description=EconDelta — daily media screen (BB indicator detection + Discord ping)
Documentation=https://github.com/clauding-lab/econdelta
After=network-online.target
Wants=network-online.target

# Retry on transient failures (LLM API / press-site brief outages).
# Up to 3 starts per 30-min window — initial fire plus 2 retries.
StartLimitIntervalSec=1800
StartLimitBurst=3

[Service]
Type=oneshot
User=adnan-local
Group=adnan-local
WorkingDirectory=/home/adnan-local/econdelta
EnvironmentFile=/etc/econdelta.env
ExecStart=/home/adnan-local/econdelta/.venv/bin/python -m scrapers.media_screen

# Resource limits to protect Notifyr and Mission Control on the same VPS
MemoryMax=500M
MemoryHigh=400M
CPUQuota=50%
TasksMax=512

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/adnan-local/econdelta/data /home/adnan-local/econdelta/logs

# LLM extraction per article adds latency. extract.py calls run_max at
# effort='high', which on the pinned box CLI maps to a high-effort Opus call —
# worst-case Opus latency under Max load can reach ~180s/call. ~5 articles ×
# (HTTP fetch + LLM call) can exceed 600s, so 600s risks a systemd kill mid-run
# after a partial insert. 1800s matches StartLimitIntervalSec (one restart fits
# the window) and mirrors the briefing service's generous LLM headroom.
TimeoutStartSec=1800
Restart=on-failure
RestartSec=300

# Logging — systemd journal + append to file
StandardOutput=append:/home/adnan-local/econdelta/logs/media_screen-systemd.log
StandardError=append:/home/adnan-local/econdelta/logs/media_screen-systemd.log
SyslogIdentifier=econdelta-media-screen

[Install]
WantedBy=multi-user.target
