# deploy/econdelta-aggregate.service
[Unit]
Description=EconDelta — Aggregate latest data across scrapers
Documentation=https://github.com/clauding-lab/econdelta
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=adnan
Group=adnan
WorkingDirectory=/home/adnan/econdelta
EnvironmentFile=/etc/econdelta.env
ExecStart=/home/adnan/econdelta/.venv/bin/python -m aggregate_latest

# Resource limits — aggregator is light
MemoryMax=250M
MemoryHigh=200M
CPUQuota=50%
TasksMax=64

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/adnan/econdelta/data /home/adnan/econdelta/logs

# Timeouts
TimeoutStartSec=30

# Logging — systemd journal + append to file
StandardOutput=append:/home/adnan/econdelta/logs/aggregate-systemd.log
StandardError=append:/home/adnan/econdelta/logs/aggregate-systemd.log
SyslogIdentifier=econdelta-aggregate

[Install]
WantedBy=multi-user.target
