# deploy/econdelta-dse.timer
[Unit]
Description=Run EconDelta DSE scraper daily at 23:11 UTC (05:11 BDT next day) — captures previous trading-day's close
Requires=econdelta-dse.service

[Timer]
OnCalendar=*-*-* 23:11:00 UTC
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
