# deploy/econdelta-aggregate-retry.timer
[Unit]
Description=EconDelta aggregate retry — fires daily at 00:10 UTC (06:10 BDT), 10 min after forex retry, so retried BB data lands in latest.json before brief.cli reads it
Requires=econdelta-aggregate.service

[Timer]
OnCalendar=*-*-* 00:10:00 UTC
Unit=econdelta-aggregate.service
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
