# deploy/econdelta-forex-retry.timer
[Unit]
Description=EconDelta forex retry — fires daily at 00:00 UTC (06:00 BDT) to backstop the primary 23:05 UTC fire (idempotent: re-fetches same-day file)
Requires=econdelta-forex.service

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Unit=econdelta-forex.service
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
