# deploy/econdelta-parse.timer
[Unit]
# Shifted from 23:15 UTC (05:15 BDT) → 22:30 UTC (04:30 BDT) to dodge the
# 05:00–06:10 BDT window where the claude CLI preflight was returning exit 1
# with empty stderr (Anthropic edge contention during US-afternoon peak).
# Tradeoff: parse runs on the prior day's fetched artifacts because fetch.timer
# still fires at 23:00 UTC; for the morning brief that's acceptable since most
# metrics (BB WSEI quarterly, NBR fiscal-year-to-date, banking NPL) refresh
# weekly or slower, and daily indicators are framed as "yesterday's close".
Description=Run EconDelta Stage 2 parse daily at 04:30 UTC (10:30 BDT)
Requires=econdelta-parse.service

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
