# deploy/econdelta-parse.timer
[Unit]
Description=Run EconDelta Stage 2 parse daily at 23:15 UTC (05:15 BDT next day)
Requires=econdelta-parse.service

[Timer]
OnCalendar=*-*-* 23:15:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
