# deploy/econdelta-fetch.timer
[Unit]
Description=Run EconDelta Stage 1 fetch daily at 00:00 UTC (06:00 BDT)
Requires=econdelta-fetch.service

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
