# deploy/econdelta-aggregate-retry.timer
[Unit]
Description=EconDelta aggregate retry — fires daily at 08:00 UTC (14:00 BDT), 1 h after primary aggregate, so any flaked aggregate at 13:00 gets a second chance before end-of-day
Requires=econdelta-aggregate.service

[Timer]
OnCalendar=*-*-* 08:00:00 UTC
Unit=econdelta-aggregate.service
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
