# deploy/econdelta-aggregate.timer
[Unit]
Description=Run EconDelta aggregator twice daily at 23:20 UTC (05:20 BDT next day) and 10:35 UTC (16:35 BDT)
Requires=econdelta-aggregate.service

[Timer]
OnCalendar=*-*-* 23:20:00 UTC
OnCalendar=*-*-* 10:35:00 UTC
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
