# deploy/econdelta-briefing.service.d/10-claude-json-writable.conf
[Service]
# The claude CLI writes ~/.claude.json each run; ProtectHome=read-only blocks it
# (EROFS) and fails the briefing. The .claude/ dir is carved out in the unit body;
# this adds the sibling .json file. See AGENT_LEARNINGS 2026-05-29.
ReadWritePaths=/home/adnan-local/.claude.json
